# discord-notifier/rss-leads-discord-notifier.timer
[Unit]
Description=Poll FreshRSS high-priority leads for Discord notifications

[Timer]
OnBootSec=10s
OnUnitActiveSec=20s
AccuracySec=1s
RandomizedDelaySec=0
Persistent=true
Unit=rss-leads-discord-notifier.service

[Install]
WantedBy=timers.target
