# discord-notifier/rss-leads-discord-notifier.service
[Unit]
Description=Notify Discord about FreshRSS high-priority leads
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/rss-leads-discord-notifier.env
StateDirectory=rss-leads-discord-notifier
StateDirectoryMode=0750
ExecStart=/usr/bin/python3 /opt/rss-leads-discord/notifier.py
