# discord-notifier/rss-leads-discord-notifier.timer
[Unit]
Description=Poll FreshRSS high-priority leads for Discord notifications

[Timer]
OnBootSec=45s
OnUnitActiveSec=60s
AccuracySec=5s
Persistent=true
Unit=rss-leads-discord-notifier.service

[Install]
WantedBy=timers.target
